# 10m.service
[Unit]
Description=10min

[Service]
ExecStart=/home/andrey/uchet/10minutes.py

