# day.service
[Unit]
Description=day

[Service]
ExecStart=/home/andrey/uchet/day.py

