# hour.service
[Unit]
Description=hour

[Service]
ExecStart=/home/andrey/uchet/hour.py

