# uchet.service
[Unit]
Description=Uchet

[Service]
ExecStart=/home/andrey/uchet/main.py

